# sickrage/DEBIAN/postinst
#!/bin/sh
# postinst script for sickrage
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

case "$1" in
    configure)
    # Make sure we give them a base ini, need to use the specified port
    if [ ! -f /data/sickrage/config/config.ini ]; then
        cp /data/sickrage/config/config.ini.dist /data/sickrage/config/config.ini
    fi

    # Chown to admin:admin
    chown -R admin:admin /apps/sickrage
    chown -R admin:admin /data/sickrage

    # Create our Config/DB Share
    rn_nml -a share:/data/sickrage

    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

###
# Enable our service
###
if [ -f /lib/systemd/system/fvapp-sickrage.service ]; then
    systemctl enable fvapp-sickrage
    systemctl start fvapp-sickrage
fi

###
# Restart Apache (required to show logo)
###
systemctl restart apache2
sleep 3

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.



exit 0
